# editions/debian/packaging/debian/ghostkernel-netd.postrm
#!/bin/sh
set -e

if [ "$1" = "remove" ]; then
    systemctl stop ghostkernel-netd.service || true
    systemctl disable ghostkernel-netd.service || true
fi

if [ "$1" = "purge" ]; then
    systemctl daemon-reload || true
fi

#DEBHELPER#

exit 0
