# editions/debian/packaging/debian/ghostkernel-auditd.postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]; then
    systemctl daemon-reload || true
    systemctl enable ghostkernel-auditd.service || true
    systemctl start ghostkernel-auditd.service || true
fi

#DEBHELPER#

exit 0
