# editions/debian/packaging/debian/ghostkernel-shardd.postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]; then
    systemctl daemon-reload || true
    systemctl enable ghostkernel-shardd.service || true
    systemctl start ghostkernel-shardd.service || true
fi

#DEBHELPER#

exit 0
