# editions/debian/packaging/debian/ghostkernel-policyd.postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]; then
    systemctl daemon-reload || true
    systemctl enable ghostkernel-policyd.service || true
    systemctl start ghostkernel-policyd.service || true
fi

#DEBHELPER#

exit 0
