# editions/debian/packaging/debian/ghostkernel-netd.postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]; then
    systemctl daemon-reload || true
    systemctl enable ghostkernel-netd.service || true
    systemctl start ghostkernel-netd.service || true
fi

#DEBHELPER#

exit 0
